# server/bin/jantac.service
[Unit]
Description=JanTac Server Module

[Service]
Type=forking
User=%%%USER%%%
Environment=DISPLAY=:0
ExecStart=/usr/bin/bash %%%PREFIX%%%/bin/start.sh
ExecStop=/usr/bin/bash %%%PREFIX%%%/bin/stop.sh
StandardOutput=null

[Install]
WantedBy=graphical.target
